# templates/squeezeboxserver.service.erb
[Unit]
Description=Logitech Media Server
After=network.target

[Service]
Type=simple
User=<%= @user_name %>
Group=<%= @user_name %>
Environment=LANG=en_US.utf8
ExecStart=<% if @custom_perl != '' %><%= @custom_perl %> <% end -%>
<%= @service_bin %>\
 --cachedir=<%= @data_dir %>\
 --prefsdir=<%= @pref_dir %>\
 --logdir=<%= @log_dir %>\
 --charset=<%= @service_charset %>\
 <%= @service_options %>
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
